# static/tools/wifi-sync/do3eet-wifi-sync.timer
[Unit]
Description=Run DO3EET WiFi Sync every hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target